# src/init/yaml_server.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec yaml_server